# com.milosz.claude-deepseek.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.milosz.claude-deepseek</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>exec "REPLACE_WITH_FULL_PATH_TO/start-litellm-deepseek.sh"</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/litellm_deepseek.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/litellm_deepseek.err.log</string>
</dict>
</plist>
